# deploy/systemd/whale-fast.timer
[Unit]
Description=Run whale-fast every 15 seconds (snapshot + diff + alert in one pass)

[Timer]
OnBootSec=1min
OnUnitActiveSec=15s
AccuracySec=2s
RandomizedDelaySec=2s
Persistent=true

[Install]
WantedBy=timers.target
